# packaging/regionlock-refresh.timer
[Unit]
Description=regionlock: periodic SDR feed snapshot refresh

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
