# packaging/regionlock.service
[Unit]
Description=regionlock: apply firewall rules from the boot snapshot
After=network-online.target nftables.service
Wants=network-online.target
# The snapshot in /etc/regionlock is self-contained; no user homedir needed.
ConditionPathExists=/etc/regionlock/config.toml

[Service]
Type=oneshot
RemainAfterExit=yes
# Runs as root: no escalation inside the unit. --system pins config and
# feed to /etc/regionlock only (a stray root homedir can never shadow the
# snapshot) and never touches the network at boot.
ExecStart=regionlock apply --yes --system
ExecStop=regionlock teardown --yes

[Install]
WantedBy=multi-user.target
